# deployment/systemd/video-service.service
# Phase 8 — run the video API as a managed background service.
# Install:
#   sudo cp deployment/systemd/video-service.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now video-service
#   sudo systemctl status video-service        # check
#   journalctl -u video-service -f             # live logs (watch for "model ready")
[Unit]
Description=Video Service (LTX / Wan)
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ai-platform/ai-inference-platform/video-service
EnvironmentFile=/opt/ai-platform/ai-inference-platform/video-service/.env
ExecStart=/opt/ai-platform/ai-inference-platform/video-service/venv/bin/uvicorn \
    app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=10
# The model loads for a few minutes on boot; don't let systemd think it hung.
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
